# profiles/gnome/aui/persistent/x86_64/upperdir/etc/systemd/system/pacman-init.service.d/override.conf
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
ConditionPathExists=!/etc/pacman.d/gnupg
