# configs/i3/airootfs/etc/NetworkManager/dispatcher.d/10-ifup
../../../../../kde/airootfs/etc/NetworkManager/dispatcher.d/10-ifup