# profiles/gnome/drive/liveusb/persistent/x86_64/upperdir/etc/systemd/system/pacman-init.service.d/override.conf
[Unit]
ConditionPathExists=!/etc/pacman.d/gnupg/openpgp-revocs.d
